# files/raspberry/monitor/supervisor/ssh.conf
[program:ssh_tunnel]
command = ssh -i /home/webcam1/.ssh/cup_desk_id_rsa -R 12345:localhost:22 cupadmin@ymmonitor.koreacentral.cloudapp.azure.com
process_name = %(program_name)s
autostart = false
autorestart = true
erxitgcodes=0
user = webcam1
stopsignal=KILL
stdout_logfile = /opt/monitor/log/ssh_tunnel.out
stderr_logfile = /opt/monitor/log/ssh_tunnel.err
